# roles/postgresql/templates/pgweb.service.j2
[Unit]
Description=pgweb 
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/pgweb
ExecStart=/usr/local/bin/pgweb  $OPTIONS

[Install]
WantedBy=multi-user.target